# mood/gunicorn-moodade.conf
start on net-device-up
stop on shutdown

respawn

setuid elspeth
chdir /home/elspeth/moodade/source

exec ../virtualenv/bin/
     --bind unix:/tmp/moodade.socket \
      moodade.wsgi:application