# static/js/moodade-gunicorn.conf
start on net-device-up
stop on shutdown

respawn

chdir /home/k220j/moodade
exec   ../virtualenv/bin/gunicorn --bind unix:/tmp/moodade.socket mysite.wsgi:application